# cron.d/zabbix_syslog_create_urls
PATH=/sbin:/usr/sbin:/bin:/usr/bin

0 1 * * *  zabbix /etc/zabbix/scripts/zabbix_syslog_create_urls.pl 2>&1
